# scripts/com.fundamental.weekly-snap.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd job: weekly score snapshot.
  Fires every Friday at 18:00 local time (IST).

  NSE Friday close is 15:30 IST; the 2.5-hour buffer lets any late-day
  price feeds settle. If the Mac is asleep at 18:00, launchd will run the
  job the next time the Mac wakes (StartCalendarInterval is catch-up
  enabled by default on macOS).

  Install:
    cp scripts/com.fundamental.weekly-snap.plist ~/Library/LaunchAgents/
    launchctl load -w ~/Library/LaunchAgents/com.fundamental.weekly-snap.plist

  See docs/SNAPSHOT_CRON.md for verify / uninstall / troubleshoot.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fundamental.weekly-snap</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/debasissahoo/Documents/Fundamental/scripts/weekly-snap.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/debasissahoo/Documents/Fundamental</string>

  <!-- Friday = 5 (Sun=0). 18:00 IST. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>5</integer>
    <key>Hour</key>
    <integer>18</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- launchd's own stdout/stderr (separate from the wrapper's per-run log). -->
  <key>StandardOutPath</key>
  <string>/Users/debasissahoo/Documents/Fundamental/logs/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/debasissahoo/Documents/Fundamental/logs/launchd.err.log</string>

  <!-- Don't fire at install time; we only want the scheduled runs. -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
